# deploy/refresh-all.timer
[Unit]
Description=India Two-Wheeler Dashboard — daily scrape refresh timer

[Timer]
# Fires every day at 04:30 UTC (~10:00 IST) with up to 1 h jitter so the
# whole world's not hammering RushLane/BikeWale at the same minute.
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
